# src/config/cron.daily/lscfg
#!/bin/bash

lscfg -a > /var/log/lscfg/lscfg.out
